# roles/ora-schema/templates/oracle-backup-ocr.service.j2
[Unit]
Description=Oracle OCR file backup
After=network.target

[Service]
Type=oneshot
ExecStart=/backup_rman/oracle_rman/scripts/ocrbck.sh

[Install]
WantedBy=multi-user.target

