# roles/ora-schema/templates/oracle-backup-archive.timer.j2
[Unit]
Description=Run oracle-backup-archive 1 times a day

[Timer]
OnCalendar=Mon-Sat *-*-* 05:00:00


[Install]
WantedBy=timers.target
