# roles/ora-schema/templates/oracle-backup-full.service.j2
[Unit]
Description=Oracle RMAN backup full
After=network.target

[Service]
User=oracle
Group=oinstall
Type=oneshot
ExecStart=/backup_rman/oracle_rman/scripts/rmanbck_multi.sh -f

[Install]
WantedBy=multi-user.target

