# roles/ora-schema/templates/oracle-backup-maintenance.service.j2
[Unit]
Description=Oracle RMAN backup maintenance
After=network.target

[Service]
User=oracle
Group=oinstall
Type=oneshot
ExecStart=/backup_rman/oracle_rman/scripts/rmanbck_multi.sh -m

[Install]
WantedBy=multi-user.target

