# roles/ora-rman/templates/oracle-backup-archive.service.j2
[Unit]
Description=Oracle RMAN backup archive log
After=network.target

[Service]
User=oracle
Group=oinstall
Type=oneshot
ExecStart=/backup_rman/oracle_rman/scripts/rmanbck_multi.sh -a

[Install]
WantedBy=multi-user.target

