# roles/ora-rman/templates/oracle-backup-incremental0.timer.j2
[Unit]
Description=Run oracle-backup-incremental0 1 times a week

[Timer]
OnCalendar=Sun *-*-* 02:00:00


[Install]
WantedBy=timers.target
